# afb-server.service
[Unit]
Description=Atomic Fishbowl Server
Documentation=https://knowledgekta.com
After=network.target mongod.service

[Service]
Type=simple
User=afb-server
WorkingDirectory=/opt/kentech/afb-server/bin
Environment="NODE_ENV=production"
Environment="SYSTEMD=1"
ExecStart=/usr/bin/node /opt/kentech/afb-server/bin/server.js
Restart=on-failure

[Install]
WantedBy=multi-user.target
